# deploy/systemd/turbodb-sync.service
# TurboDB CDC sync daemon (PostgreSQL -> engine).
# Install: cp turbodb-sync.service /etc/systemd/system/
#          cp sync.env.example /etc/turbodb/sync.env   # holds TURBODB_PG_DSN
#          install -m 0640 sync.yaml /etc/turbodb/sync.yaml
#          useradd --system --home /var/lib/turbodb-sync --shell /usr/sbin/nologin turbodb-sync
#          systemctl daemon-reload && systemctl enable --now turbodb-sync
[Unit]
Description=TurboDB CDC replication sync
Documentation=https://github.com/gabrielleeyj/turbodb
After=network-online.target turbodb-engine.service
Wants=network-online.target

[Service]
Type=simple
User=turbodb-sync
Group=turbodb-sync
# sync.env carries TURBODB_PG_DSN. The DSN must live in the environment
# file (mode 0640), never on the command line where it is visible in ps.
EnvironmentFile=/etc/turbodb/sync.env
ExecStartPre=/usr/local/bin/turbodb-sync check-config --config /etc/turbodb/sync.yaml
ExecStart=/usr/local/bin/turbodb-sync run \
    --config /etc/turbodb/sync.yaml \
    --slot turbodb_sync \
    --publication turbodb_pub \
    --engine 127.0.0.1:7080 \
    --checkpoint /var/lib/turbodb-sync/sync.ckpt
# Sync exits by design when the engine becomes unreachable (circuit
# breaker); unconditional restart with a short delay is the supervision
# contract, mirroring the soak harness.
Restart=always
RestartSec=2
StateDirectory=turbodb-sync
LimitNOFILE=65536

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/turbodb-sync
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes

[Install]
WantedBy=multi-user.target
